# schedule2.cron
56 05 * * *   /scratch/crsusr/tools/allinone.sh & >> /tmp/log.txt
